# kafka-generator.ini
[program:hosts-gen]
command=/usr/bin/python /root/data_gen/src.py 127.0.0.1:9092 100000 2 2
stdout_logfile=/var/log/hosts.log
stderr_logfile=/var/log/hosts_err.log
umask=022
autostart=true
autorestart=true
startsecs=10
startretries=3